# simple-blog.service
[Unit]
Description=Simple Blog Flask Application
Documentation=https://github.com/lbxxgn/my-blog
After=network.target

[Service]
Type=simple
# 建议使用专用用户运行（如果使用 root，请保持原样）
User=root
Group=root

# 项目目录（根据实际部署路径修改）
WorkingDirectory=/root/my-blog

# 环境变量
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="FLASK_APP=app.py"
Environment="FLASK_ENV=production"
Environment="PYTHONPATH=/root/my-blog"
Environment="PORT=5001"

# 从 .env 文件加载环境变量（不存在不报错）
# 注意：- 前缀表示文件可选
EnvironmentFile=-/root/my-blog/.env

# 启动命令
ExecStart=/usr/bin/python3 /root/my-blog/backend/app.py

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# 日志配置（相对于 WorkingDirectory）
StandardOutput=append:/root/my-blog/logs/app.log
StandardError=append:/root/my-blog/logs/error.log

# 安全设置（生产环境建议启用）
# NoNewPrivileges=true
# PrivateTmp=true

# 资源限制
LimitNOFILE=65536
# MemoryLimit=512M

[Install]
WantedBy=multi-user.target
